# CNLdisplay.service
[Unit]
Description=CNL Display and Alarm service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/pi/
ExecStart=/usr/bin/python2.7 -m decoding-contour-next-link.read_minimed_next24
# Below will make service stop in 9hrs (7min commented out  below that for testing)
RuntimeMaxSec=32400
# RuntimeMaxSec=420
ExecStopPost=/home/pi/decoding-contour-next-link/displayERR.py

[Install]
WantedBy=multi-user.target
